# templates/chezmoi-update@.timer.j2
[Unit]
Description=Update dotfiles for user %I

[Timer]
OnBootSec=1min
OnCalendar={{ chezmoi_update_cron }}
Persistent=false
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
